# deploy/com.vickers.mcp4chatgpt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.vickers.mcp4chatgpt</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/Cellar/python@3.14/3.14.5/Frameworks/Python.framework/Versions/3.14/bin/python3.14</string>
		<string>-m</string>
		<string>mcp4chatgpt.server</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>/Users/vickers/Documents/MCP_Creator/MCP4ChatGPT/src</string>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/vickers/Documents/MCP_Creator/MCP4ChatGPT</string>
	<key>RunAtLoad</key>
	<false/>
	<key>KeepAlive</key>
	<false/>
	<key>StandardOutPath</key>
	<string>/Users/vickers/Documents/MCP_Creator/MCP4ChatGPT/logs/service.out.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/vickers/Documents/MCP_Creator/MCP4ChatGPT/logs/service.err.log</string>
</dict>
</plist>
